# database-server-centos/build/supervisor-mysqld.conf

[program:mysqld]
command=/usr/sbin/mysqld