# build/cloudstack-supervisor.conf.j2
[program:cloudstack]
command=/home/cloudstack/run-cloudstack.sh
username=cloudstack
loglevel=debug
directory=/home/cloudstack/
stdout_events_enabled = true